# deploy/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=/usr/sbin/php5-fpm

[program:nginx]
command=nginx

[program:nginx]
command=nginx

[program:mysqld]
command=mysqld
