# deploy/runner/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:teddy-runner]
directory=/opt/teddy/runner/bin
command=java -jar teddy.jar
